# deploy/systemd/wg-hope-monitor.service
[Unit]
Description=WG Hope Monitor Web UI
After=network-online.target wg-quick@wg0.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/wg-hope-bot
EnvironmentFile=/opt/wg-hope-bot/.env
ExecStart=/bin/bash -lc '/opt/wg-hope-bot/.venv/bin/python -m uvicorn monitor.app:app --host 0.0.0.0 --port ${WEB_PORT:-8080}'
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
